# lpcsnoop.service.in
[Unit]
Description=LPC Snoop Daemon
After=systemd-modules-load.service

[Service]
Restart=always
ExecStart=@bindir@/snoopd -b 1 -h 4

[Install]
WantedBy=@SYSTEMD_TARGET@
